# systemd/orion-dns-ha-health.service
# =============================================================================
# Orion DNS HA - Health Check Service
# =============================================================================
#
# This systemd service runs the host-level DNS health check script.
# It should be triggered by the companion timer (orion-dns-ha-health.timer).
#
# Installation:
#   1. Copy this file to /etc/systemd/system/orion-dns-ha-health.service
#   2. Copy orion-dns-ha-health.timer to /etc/systemd/system/
#   3. Run: sudo systemctl daemon-reload
#   4. Run: sudo systemctl enable --now orion-dns-ha-health.timer
#
# Configuration:
#   Edit the Environment lines below or override via /etc/default/orion-dns-ha
#
# =============================================================================

[Unit]
Description=Orion DNS HA - health check and auto-heal
Documentation=https://github.com/orionsentinel/Orion-sentinel-ha-dns
Wants=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/orion-dns-ha
ExecStart=/opt/orion-dns-ha/ops/orion-dns-health.sh
TimeoutStartSec=60

# Environment configuration (can be overridden)
Environment="REPO_DIR=/opt/orion-dns-ha"
Environment="PIHOLE_CONTAINER_NAME=pihole_unbound"
Environment="KEEPALIVED_CONTAINER_NAME=keepalived"
Environment="HEALTH_FAIL_THRESHOLD=2"

# Optional: Load from environment file if it exists
EnvironmentFile=-/etc/default/orion-dns-ha

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=orion-dns-health

[Install]
WantedBy=multi-user.target
