# systemd/orion-dns-ha-backup-node.service
# =============================================================================
# Orion DNS HA - Backup Node Service
# =============================================================================
#
# This systemd service manages the Docker Compose stack on the BACKUP node.
# Use this on the node that should be BACKUP by default (lower priority).
#
# Installation (on BACKUP node only):
#   1. Copy this file to /etc/systemd/system/orion-dns-ha-backup.service
#   2. Run: sudo systemctl daemon-reload
#   3. Run: sudo systemctl enable --now orion-dns-ha-backup.service
#
# NOTE: Only install ONE of these services per node:
#   - orion-dns-ha-primary.service - for the primary/MASTER node
#   - orion-dns-ha-backup.service (this file) - for the secondary/BACKUP node
#
# Verification:
#   sudo systemctl status orion-dns-ha-backup.service
#   docker ps
#
# =============================================================================

[Unit]
Description=Orion DNS HA Stack - Backup Node
Documentation=https://github.com/orionsentinel/Orion-sentinel-ha-dns
Requires=docker.service
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/orion-dns-ha

# Start the stack with the backup profile
ExecStart=/usr/bin/docker compose --profile two-node-ha-backup up -d

# Stop the stack gracefully
ExecStop=/usr/bin/docker compose --profile two-node-ha-backup down

# Give containers time to start/stop
TimeoutStartSec=120
TimeoutStopSec=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=orion-dns-ha

[Install]
WantedBy=multi-user.target
